# assets/text/dns_traffic.service
[Unit]
Description=tcpdump DNS traffic

[Service]
Type=simple
ExecStart=/sbin/dns_traffic.bsh

[Install]
WantedBy=multi-user.target
